# setup/systemd/bushidan-evolution.service
[Unit]
Description=武士団 週次スキル進化バッチ
After=network.target postgresql.service

[Service]
Type=oneshot
User=claude
WorkingDirectory=/mnt/Bushidan-Multi-Agent
ExecStart=/mnt/Bushidan-Multi-Agent/.venv/bin/python /mnt/Bushidan-Multi-Agent/scripts/weekly_evolution.py
EnvironmentFile=/mnt/Bushidan-Multi-Agent/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bushidan-evolution

[Install]
WantedBy=multi-user.target
